# deployments/memory-api.service
[Unit]
Description=Spatial Query REST API Server
After=network.target

[Service]
Type=simple
User=
Group=
WorkingDirectory=
Environment=PYTHONPATH=
Environment=VIRTUAL_ENV=
Environment=PATH=
EnvironmentFile=
ExecStart=
Restart=
RestartSec=3
StandardOutput=
StandardError=

[Install]
WantedBy=multi-user.target 